# install/usr/local/lib/systemd/system/dockside.service
[Unit]
Description=Dockside container stack
Requires=docker.service
After=docker.service

[Service]
Type=notify
NotifyAccess=main
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300
ExecStart=/usr/local/lib/dockside/dockside-network-firewall.sh --daemon
ExecStartPost=/usr/bin/docker compose -f /usr/local/lib/dockside/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f /usr/local/lib/dockside/docker-compose.yml down
ExecReload=/usr/local/lib/dockside/dockside-network-firewall.sh --refresh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
